# docker/tasks/cron/cron.d
*/5 * * * * dotnet /app/build/Tasks.dll run ImportProducts
*/5 * * * * dotnet /app/build/Tasks.dll run ImportOrder